# src/frc_payload_launcher/launch_files/fastdds-shm.service
[Unit]
Description=Configure FastDDS shared-memory transport profile
DefaultDependencies=no
Before=sysinit.target
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'mkdir -p /etc/environment.d && \
    printf "FASTRTPS_DEFAULT_PROFILES_FILE=/etc/ros/fastdds_shm.xml\n" \
    > /etc/environment.d/fastdds-shm.conf'

[Install]
WantedBy=sysinit.target
